# secretary-frontend.service
[Unit]
Description=Secretary Partner AI Frontend (Vite)
After=network.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/Secretary_Partner_AI/frontend
Environment="PATH=/home/ec2-user/.nvm/versions/node/v20.18.3/bin:/usr/local/bin:/usr/bin:/bin"
Environment="NODE_ENV=production"
ExecStart=/home/ec2-user/.nvm/versions/node/v20.18.3/bin/npm run preview -- --host 0.0.0.0 --port 5173
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
